# .config/systemd/user/niri-autostart.service
[Service]
ExecStart=/usr/bin/niri-session

[Install]
WantedBy=default.target
