# .config/systemd/user/niri-autostart.service
[Service]
ExecStart=prime-run /usr/bin/niri-session

[Install]
WantedBy=default.target
